# ts/ops/com.whitemountain.check-waitlists.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.whitemountain.check-waitlists</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/notronwest/data/web/wmpc/projects/court-reserve-scheduler/ts/ops/run-check-waitlists.sh</string>
    </array>

    <!-- Daily at 9, 11, 13, 15, 17. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/notronwest/Library/Logs/court_reserve/launchd_check_waitlists.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/notronwest/Library/Logs/court_reserve/launchd_check_waitlists_err.log</string>
</dict>
</plist>
